# config/agent-governance.timer
[Unit]
Description=Daily AI Agent Governance Sync
Documentation=https://github.com/user/agent-governance/blob/main/docs/PHASES.md

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=agent-governance.service

[Install]
WantedBy=timers.target
